# vow.service
[Unit]
Description=VPN on wifi monitor
After=network.target

[Service]
ExecStart=/usr/local/bin/vow

[Install]
WantedBy=multi-user.target
